# files/elements/infra-cloud-bridge/static/etc/init/create_bridge.conf
description "create bridge for infra cloud"

console output

start on stopped glean

script
    exec python /opt/create_bridge.py
end script
